# config/supervisord.d/bn-service-queue-activity.ini
[program:service-queue-activity]
user = www
numprocs=2
process_name=%(program_name)s_%(process_num)02d
command = /usr/bin/php  /www/wwwroot/jihui_group/artisan queue:work redis --sleep=3 --tries=3 --queue=HandleSecKill,CloseSecKillTrade,CloseTrade,GroupClearUser,SplitTrade,DisReward
autostart = true
startsecs = 5
autorestart = true
startretries = 3
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 10
stdout_logfile = /var/log/hym_queue.log